# etc/s6-rc.d/cron/run
#!/command/with-contenv bash

set -e

cd /app

exec supercronic etc/crontab
