# etc/crontab
@weekly ./manage.py update_index
@daily ./manage.py clearsessions
@daily ./manage.py update_unsplash_photos
*/10 * * * * ./manage.py publish_scheduled_pages
