# etc/s6-rc.d/nginx/run
#!/command/with-contenv bash

set -e

cd /app

exec nginx -g "daemon off;"
