# etc/s6-rc.d/rq/run
#!/command/with-contenv bash

set -e

cd /app

exec python manage.py db_worker -v3 --interval 10
