# etc/s6-rc.d/rq/run
#!/command/with-contenv bash

set -e

cd /app

exec s6-setuidgid website python manage.py rqworker --with-scheduler
