# etc/s6-rc.d/migrate/up
s6-setuidgid website with-contenv bash -c "cd /app && python manage.py migrate --noinput"
